# supervisord.conf
{% if settings.ENVIRONMENT == 'dev' %}

[program:runserver]
command={{ PROJECT_DIR }}/run.sh {{ PROJECT_DIR }}/manage.py runserver 0.0.0.0:5000
stopsignal=KILL
stopasgroup=true

{% else %}

[program:gunicorn]
command={{ PROJECT_DIR }}/run.sh {{ PROJECT_DIR }}/manage.py run_gunicorn --workers=2 --access-logfile=gunicorn_access.log -p {{ PROJECT_DIR }}/gunicorn.pid
stopsignal=KILL
stopasgroup=true

{% endif %}
